# docker/api.supervisor.conf
[program:apisense-api]
command = apisense api
numprocs = 1
autostart = true
autorestart = true