# docker/api.supervisor.conf
[program:odm-api]
command = odh-data-monitor api
numprocs = 1
autostart = true
autorestart = true