# docker/daemon.supervisor.conf
[program:apisense-daemon]
command = apisense daemon start
numprocs = 1
autostart = true
autorestart = true