# docker/ssh.supervisor.conf
[program:odm-ssh]
command = odh-data-monitor ssh
numprocs = 1
autostart = true
autorestart = true