# docker/ssh.supervisor.conf
[program:apisense-ssh]
command = apisense ssh
numprocs = 1
autostart = true
autorestart = true